# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
directory = /opt/swagmail/git
command = /opt/swagmail/env/bin/uwsgi --ini=/opt/swagmail/git/wsgi.ini
